# scripts/systemd/maschina-update.timer
[Unit]
Description=Maschina auto-update timer (every 5 minutes)
After=network-online.target

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
